# files/mysqlbix.service
[Unit]
Description=
After=sshd.service
Requires=sshd.service
 
[Service]
Type=simple
ExecStart=
ExecStop=
User=<enter the user you want this to run as>
Group=<enter the group you want this to run as>
 
[Install]
WantedBy=multi-user.target